# conf/redis.openbroadcast.ch.supervised.conf
[program:redis.openbroadcast.ch]
directory=/root/
command=/usr/local/bin/redis-server /etc/redis/redis.conf
user=root
autostart=true
autorestart=true
redirect_stderr=True
environment=HOME='/root/'