# website/supervisord.conf
[program:webserver]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py runserver 0.0.0.0:8080

[program:celeryworker]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py celeryd -c 3 -Q celery
